# .crontab
MAILTO=""

# ---------------------
# command variable division
# ---------------------

PHP="/usr/bin/php"

# ---------------------
# wiz systems
# ---------------------

HANDLER="/home/wiz/g/application/wiz/public/index.php"

# migrate redmine status and strip assignment
*/1 * * * * ${PHP} ${HANDLER} cron convert_status_and_strip_assignment 
