# .crontab
MAILTO=""

# ---------------------
# command variable division
# ---------------------

PHP="/usr/bin/php"

# ---------------------
# wiz systems
# ---------------------

HANDLER="/home/wiz/g/application/wiz/public/index.php"

# feed addup data
30 */3 * * * ${PHP} ${HANDLER} cron feed_addup_data
